# scripts/crontab
#!/bin/bash
#PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
#@reboot /root/.config/xmodmap_vi.sh
* * * * * ~/.config/tnt/bat.sh
