# scripts/crontab
#!/bin/sh

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root export DISPLAY=:0 && /bin/sh /root/.config/bat.sh
