# mod/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --nohostname --nohints --noissue --skip-login --login-options "-f root" %I 38400 linux
